# ansible/roles/adafruit-unifi/templates/adafruit-unifi-supervisor.conf
[program:adafruit_unifi]
command='{{ git_path }}/adafruit-unifi/main.py'
autostart=true
autorestart=true
